# src/sup-do.timer
[Unit]
Description=sup-do: run keepalive at 04/10/16/22, recovering missed runs

[Timer]
OnCalendar=*-*-* 04,10,16,22:00:00
Persistent=true

[Install]
WantedBy=timers.target