# src/supabase_keepalive.service
[Unit]
Description=Supabase free tier keepalive logger

[Service]
Type=oneshot
EnvironmentFile=%h/.config/supabase_keepalive.env
ExecStart=%h/.local/bin/supabase_keepalive.sh